# celery-worker.service
[Unit]
Description=Celery Worker for WoW Guild Analytics
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=guildmaestro
Group=guildmaestro
WorkingDirectory=/var/www/guildmaestro
Environment="PATH=/var/www/guildmaestro/venv/bin"
Environment="PYTHONPATH=/var/www/guildmaestro"

# Load environment variables from .env file
EnvironmentFile=/var/www/guildmaestro/.env

# Celery worker command
# Note: Increase --concurrency based on your database backend:
# - SQLite: Use 1-2 (to avoid database locks)
# - PostgreSQL: Use 4+ (can handle concurrent connections)
ExecStart=/var/www/guildmaestro/venv/bin/celery -A app.celery_config.celery worker \
    --loglevel=info \
    --logfile=/var/www/guildmaestro/logs/celery-worker.log \
    --pidfile=/var/run/celery/celery-worker.pid \
    --concurrency=4 \
    --max-tasks-per-child=1000

# Runtime directory for PID file
RuntimeDirectory=celery
RuntimeDirectoryMode=0755

# Process management
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=60

# Security settings
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=append:/var/www/guildmaestro/logs/celery-worker.log
StandardError=append:/var/www/guildmaestro/logs/celery-worker-error.log

[Install]
WantedBy=multi-user.target
