# celery-beat.service
[Unit]
Description=Celery Beat Scheduler for WoW Guild Analytics
After=network.target redis-server.service celery-worker.service
Requires=redis-server.service

[Service]
Type=simple
User=guildmaestro
Group=guildmaestro
WorkingDirectory=/var/www/guildmaestro
Environment="PATH=/var/www/guildmaestro/venv/bin"
Environment="PYTHONPATH=/var/www/guildmaestro"

# Load environment variables from .env file
EnvironmentFile=/var/www/guildmaestro/.env

# Celery Beat command
ExecStart=/var/www/guildmaestro/venv/bin/celery -A app.celery_config.celery beat \
    --loglevel=info \
    --logfile=/var/www/guildmaestro/logs/celery-beat.log \
    --pidfile=/var/run/celery/celery-beat.pid \
    --schedule=/var/run/celery/celerybeat-schedule.db

# Runtime directory for PID file and schedule DB
RuntimeDirectory=celery
RuntimeDirectoryMode=0755

# Process management
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=append:/var/www/guildmaestro/logs/celery-beat.log
StandardError=append:/var/www/guildmaestro/logs/celery-beat-error.log

[Install]
WantedBy=multi-user.target
