# flower.service
[Unit]
Description=Flower - Celery Monitoring Dashboard for WoW Guild Analytics
After=network.target redis-server.service celery-worker.service
Requires=redis-server.service

[Service]
Type=simple
User=guildmaestro
Group=guildmaestro
WorkingDirectory=/var/www/guildmaestro
Environment="PATH=/var/www/guildmaestro/venv/bin"
Environment="PYTHONPATH=/var/www/guildmaestro"

# Load environment variables from .env file
EnvironmentFile=/var/www/guildmaestro/.env

# Flower command (monitoring dashboard on port 5555)
ExecStart=/var/www/guildmaestro/venv/bin/celery -A app.celery_config.celery flower \
    --port=5555 \
    --url_prefix=flower \
    --basic_auth=admin:${FLOWER_PASSWORD}

# Process management
Restart=always
RestartSec=10
KillSignal=SIGTERM

# Security settings
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=append:/var/www/guildmaestro/logs/flower.log
StandardError=append:/var/www/guildmaestro/logs/flower-error.log

[Install]
WantedBy=multi-user.target
